# conf/gatilho.service
[Service]
ExecStart=/bin/bash /opt/bravo-prepare/percursor.sh

[Install]
WantedBy=default.target